# data/etc/crontabs/root
0 0 * * * /etc/init.d/acme start
* * * * * /etc/scripts/collect-metrics.sh
